# clio-vigil/config/systemd/clio-vigil-uap.timer
[Unit]
Description=Clio Vigil UAP — kör dagligen 06:30

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
